# server/SHserver.service
[Unit]
Description=Smart House Server

[Service]
Type=forking
ExecStart=/SHserver
Restart=on-failure

[Install]
WantedBy=multi-user.target
